# telegram-bot.service
[Unit]
Description=Telegram File Sharing Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/telegram-file-bot
Environment=BOT_TOKEN=your_bot_token
Environment=BASE_URL=http://your-server-ip:8080
Environment=FILES_DIR=/home/ubuntu/telegram-file-bot/files
ExecStart=/usr/bin/python3 /home/ubuntu/telegram-file-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 